# roles/wake_on_lan/files/etc/systemd/system/power_off.timer
[Unit]
Description=Timer to power off devices

[Timer]
OnCalendar=1:45
Persistent=true

[Install]
WantedBy=timers.target
