# roles/wake_on_lan/files/etc/systemd/system/power_on.timer
[Unit]
Description=Timer to power on devices

[Timer]
OnCalendar=02:30
Persistent=true

[Install]
WantedBy=timers.target
