# roles/server1_services/templates/systemd/power_off.service.j2
[Unit]
Description=Power Off Devices
After=network.target

[Service]
ExecStart=/init/power_off.sh
User=root
Restart=on-failure
RestartSec=5s
StartLimitInterval=30s
StartLimitBurst=5
